# vcserver.service.in
[Unit]
Description=Visual Concepts Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/vcserver SYSCONFDIR/vcserver.cfg
StandardOutput=append:/var/log/vcserver.log

[Install]
WantedBy=multi-user.target
